# etc/periodic/daily/backup
#!/bin/sh

# Important! Script have to be without extension!

/usr/bin/coverme backup -c /etc/coverme/config.yml
